# root/etc/crontabs/abc
# min   hour    day     month   weekday command
2       3       *       *       *       /app/update.sh 2>&1
